# server/deploy/digitalocean/kolkhoz-greenfield.service
[Unit]
Description=Kolkhoz online server
After=network-online.target kolkhoz-greenfield-redis.service
Requires=kolkhoz-greenfield-redis.service
Wants=network-online.target

[Service]
Type=simple
User=kolkhoz-greenfield
Group=kolkhoz-greenfield
WorkingDirectory=/opt/kolkhoz-greenfield
EnvironmentFile=/etc/kolkhoz-greenfield.env
EnvironmentFile=/etc/kolkhoz-greenfield-build.env
ExecStart=/opt/kolkhoz-greenfield/.venv/bin/uvicorn server.kolkhoz_server.production:create_asgi_application --factory --host 127.0.0.1 --port 18787 --workers 1 --ws-ping-interval 20 --ws-ping-timeout 20
Restart=on-failure
RestartSec=5
TimeoutStopSec=20
CPUQuota=100%
MemoryMax=300M
MemoryHigh=260M
TasksMax=128
LimitNOFILE=8192
OOMPolicy=stop
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=full
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictSUIDSGID=true
ReadWritePaths=/opt/kolkhoz-greenfield/research/.build

[Install]
WantedBy=multi-user.target
