# server/deploy/digitalocean/kolkhoz-greenfield-redis.service
[Unit]
Description=Loopback Redis for Kolkhoz online server
After=network.target

[Service]
Type=notify
User=redis
Group=redis
ExecStart=/usr/bin/redis-server /etc/redis/kolkhoz-greenfield.conf --supervised systemd
ExecStartPost=/usr/bin/redis-cli -p 16379 CONFIG SET maxclients 1000
ExecStop=/usr/bin/redis-cli -p 16379 shutdown
Restart=on-failure
CPUQuota=10%
MemoryMax=96M
TasksMax=64
LimitNOFILE=4096
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=full
ReadWritePaths=/var/lib/redis-greenfield /var/log/redis

[Install]
WantedBy=multi-user.target
