# server/deploy/digitalocean/kolkhoz-ai-canary.service
[Unit]
Description=Kolkhoz real C-engine AI policy canary
After=kolkhoz-greenfield.service

[Service]
Type=oneshot
User=kolkhoz-greenfield
Group=kolkhoz-greenfield
WorkingDirectory=/opt/kolkhoz-greenfield
EnvironmentFile=/etc/kolkhoz-greenfield.env
ExecStart=/usr/local/sbin/kolkhoz-ai-canary
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=full
ReadWritePaths=/opt/kolkhoz-greenfield/research/.build
