# server/deploy/digitalocean/kolkhoz-health-watch.service
[Unit]
Description=Kolkhoz production canary and readiness alarm
After=kolkhoz-greenfield.service

[Service]
Type=oneshot
EnvironmentFile=/etc/kolkhoz-greenfield.env
ExecStart=/usr/local/sbin/kolkhoz-health-watch
