# server/deploy/digitalocean/kolkhoz-unconfirmed-account-cleanup.service
[Unit]
Description=Delete stale unconfirmed Kolkhoz Supabase accounts
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=kolkhoz-server
Group=kolkhoz-server
WorkingDirectory=/opt/kolkhoz-server
EnvironmentFile=/etc/kolkhoz-server.env
ExecStart=/opt/kolkhoz-server/.venv/bin/python -m server.tools.cleanup_unconfirmed_accounts --older-than-days 7 --delete
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
