# server/deploy/digitalocean/kolkhoz-ai-canary.service
[Unit]
Description=Kolkhoz real C-engine AI policy canary
After=kolkhoz-server.service

[Service]
Type=oneshot
User=kolkhoz-server
Group=kolkhoz-server
WorkingDirectory=/opt/kolkhoz-server
EnvironmentFile=/etc/kolkhoz-server.env
ExecStart=/usr/local/sbin/kolkhoz-ai-canary
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=full
ReadWritePaths=/opt/kolkhoz-server/research/.build
