# deploy/digitalocean/kolkhoz-online.service
[Unit]
Description=Kolkhoz online session server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=kolkhoz
Group=kolkhoz
WorkingDirectory=/opt/kolkhoz
EnvironmentFile=/etc/kolkhoz-online.env
ExecStart=/opt/kolkhoz/.venv/bin/python -m research.kolkhoz_research.cli serve-online --host 127.0.0.1 --port 8787 --rebuild
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ReadWritePaths=/opt/kolkhoz

[Install]
WantedBy=multi-user.target
